# roles/mailinabox/files/etc/profile.d/miab-env.sh
#!/usr/bin/env bash
export NONINTERACTIVE=1
export PRIMARY_HOSTNAME=box.yourdomain.tld
