# roles/common/files/etc/cron.monthly/hadori.sh
#!/bin/sh
# Hardlink specified directories using hadori
/usr/bin/hadori --verbose /bin /sbin /lib /lib64 /usr /opt
