# roles/common/files/etc/cron.monthly/hadori
#!/bin/sh
/usr/bin/hadori --verbose /bin /sbin /lib /lib64 /usr /opt
