# roles/common/files/etc/cron.daily/apt-get_update
#!/bin/sh
/usr/bin/apt-get update --quiet=2
